# qr3_81.conf
[program:qr3_81]
command=python3 /home/curator3/qr3_81/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_81/test.err.log
stdout_logfile=/home/curator3/qr3_81/test.err.log
